# rootfs/etc/services.d/unitd/run
#!/usr/bin/env sh

exec /opt/unit/sbin/unitd --no-daemon --state "/var/lib/unit/state" --modules "/opt/unit/modules" --control "unix:${UNIT_SOCKET}"
