# rootfs/etc/s6-overlay/s6-rc.d/init-unitd-load-secrets/up
/etc/s6-overlay/s6-rc.d/init-unitd-load-secrets/run